# my_msi_driver.service.in
[Unit]
Description=MSI MEG CoreLiquid S360 Daemon

[Service]
Type=simple
ExecStart=@INSTALL_PREFIX@/bin/my_msi_driver -M 5 startd
StandardOutput=journal
StandardError=journal
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
KillMode=mixed
Restart=on-failure

[Install]
WantedBy=multi-user.target
